# caravel-apps.service
[Unit]
Description=Caravel Clinic ZIP Code Applications
After=network.target

[Service]
Type=simple
User=viriditytech
WorkingDirectory=/home/viriditytech/apps/zip-optimizer
Environment=PATH=/home/viriditytech/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 start_azure_apps.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 